# rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/execlineb -P
/usr/sbin/nginx -c /config/nginx/nginx.conf
