# rootfs/etc/s6-overlay/s6-rc.d/spawn-fcgi/run
#!/command/execlineb -P
/usr/bin/spawn-fcgi -u yaapc -g yaapc -s /var/run/fcgiwrap.sock -n -- /usr/bin/fcgiwrap
