# rootfs/etc/services.d/spawn-fcgi/run
#!/usr/bin/with-contenv bash

exec spawn-fcgi -u yaapc -g yaapc -s /var/run/fcgiwrap.sock -n -- /usr/bin/fcgiwrap